# root/etc/services.d/jupyterlab/run
#!/usr/bin/with-contenv bash

exec jupyter-lab --ip 0.0.0.0 --no-browser --autoreload --allow-root --notebook-dir /home/jovyan/workspace
